# private_dot_config/etc-backup/files/init.d/executable_AmneziaVPNd
#!/sbin/openrc-run

name="AmneziaVPNd daemon"
description="AmneziaVPN Service"
pidfile="/run/amnezia-vpn.pid"

depend() {
    # netifrc
    # after net
    after NetworkManager
}

start() {
    ebegin "Starting AmneziaVPN Service"
    export LD_LIBRARY_PATH=/opt/AmneziaVPN/client/lib
    start-stop-daemon --start --background --make-pidfile --pidfile "$pidfile" --exec /opt/AmneziaVPN/service/AmneziaVPN-service.sh
    eend $?
}

stop() {
    ebegin "Stopping AmneziaVPN Service"
    start-stop-daemon --stop --pidfile "$pidfile"
    eend $?
}
